# infrastructure/systemd/thesis-vllm-metrics-proxy-heavy.service
[Unit]
Description=Thesis vLLM Metrics Proxy (Heavy Host)
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=simple
User=dogukan
WorkingDirectory=/home/dogukan/thesis-proxy-src
Environment=PYTHONPATH=/home/dogukan/thesis-proxy-src
Environment=THESIS_VLLM_UPSTREAM=http://127.0.0.1:8001/v1
Environment=THESIS_PROXY_PORT=8000
ExecStart=/home/dogukan/thesis-proxy-venv/bin/python /home/dogukan/thesis-proxy-src/infrastructure/vllm/metrics_proxy.py
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
